# coap-server.service
[Unit]
Description=CoAP Server
After=network.target

[Service]
Type=simple
User=coap
Group=coap
WorkingDirectory=/opt/coap-server
ExecStart=/usr/local/bin/bun run src/index.ts
Restart=on-failure
RestartSec=5

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
